# Ansible/backup_mediawiki_data_service/templates/backup_media_dir.timer.j2
[Unit]
Description=Run Mediawiki  Backup Script every 5 minutes

[Timer]
OnBootSec=7min
OnUnitActiveSec=7min
Persistent=true

[Install]
WantedBy=timers.target